# launch_server.service
[Unit]
After=network.target

[Service]
ExecStart=/home/phoenix/phoenix-server/PhoenixServer.sh
Restart=on-failure

[Install]
WantedBy=default.target